# [project start]/jesmyl_soki.service
[Unit]
Description=The Soki Service

[Service]
Restart=on-failure
RestartSec=5s
ExecStart=node /var/www/jesmyl.ru/back.index.js