# [project start]/jesmyl_soki.service
[Unit]
Description=The Soki Service

[Service]
Restart=on-failure
RestartSec=5s
ExecStart=ts-node /var/www/{DOMAIN}/complect/soki/SokiServer.ts