# [project start]/jmsg.service
[Unit]
Description=The Msg Soki Service

[Service]
Restart=on-failure
RestartSec=5s
ExecStart=node /var/www/html/back.index.js

#3359